# supervisord.conf
[program:collector]
command=python runcollection.py
autorestart=true
